# systemd/user/session.target
[Unit]
Description=The actual session
AllowIsolate=true
